# files/etc/sv/nginx/run
#!/bin/bash
exec /usr/sbin/nginx -g "daemon off;"
