# files/etc/sv/php/run
#!/bin/sh

mkdir -p /run/php/
exec /usr/sbin/php-fpm7.0 --nodaemonize -R

